# config/autostart/dunst.desktop
[Desktop Entry]
Encoding=UTF-8
Type=Application
Name=dunst
Exec=dunst
StartupNotify=false
Terminal=false
Hidden=false
